# private_dot_config/private_systemd/private_user/private_bitwarden-ssh-agent.socket
[Unit]
Description=Bitwarden SSH agent socket

[Socket]
ListenStream=%t/bitwarden-ssh-agent.sock
SocketMode=0600
Accept=yes
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
