# Library/LaunchAgents/com.ollama.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.ollama.serve</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Users/anh.pham/.pixi/bin/ollama</string>
            <string>serve</string>
        </array>

        <key>EnvironmentVariables</key>
        <dict>
            <key>OLLAMA_FLASH_ATTENTION</key>
            <string>true</string>
            <key>OLLAMA_KEEP_ALIVE</key>
            <string>30m</string>
            <key>OLLAMA_CONTEXT_LENGTH</key>
            <string>8192</string>
            <key>OLLAMA_NEW_ENGINE</key>
            <string>true</string>
        </dict>

        <key>RunAtLoad</key>
        <true />
        <key>KeepAlive</key>
        <dict>
            <key>SuccessfulExit</key>
            <false />
        </dict>
    </dict>
</plist>